# deploy/cute-relay.service
[Unit]
Description=cute pixel relay (WebSocket transit, console <-> cute_pixel app)
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/cute/relay
ExecStart=/usr/bin/node dist/server.cjs
Restart=on-failure
RestartSec=3s
Environment=NODE_ENV=production
Environment=PORT=8080
Environment=HOST=127.0.0.1

# 安全 hardening:relay 只读自己目录,无 sudo,无新权限。
# 注:不能用 ProtectHome=true,因为 WorkingDirectory 在 /home/ubuntu/,会 CHDIR 失败。
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/home/ubuntu/cute/relay
PrivateTmp=true

[Install]
WantedBy=multi-user.target
